# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

; On Render, MCP_DIRECT=true means tools are imported directly
; into the API process — no separate MCP server processes needed.
; This saves ~300MB RAM and eliminates all port/timing issues.

[program:api]
command=python -u -m uvicorn src.main:app --host 0.0.0.0 --port 10000
directory=/app
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED=1
